# crontab
0 */1 * * * /usr/local/bin/python /parsing/parser.py >> /var/log/cron.log 2>&1
